# docker/app/scheduler/supervisord.conf
[program:schedule]
command=/usr/local/etc/scheduler/run.sh
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
